# net/portmap/init.d/portmap
#!/bin/sh
#
# portmap
#
# chkconfig: 345 11 89
# description: the portmapper is the RPC port broker
# processname: /sbin/portmap
#

. /lib/lsb/init-functions

