# doc-tools/man/cron.weekly/makewhatis
#!/bin/bash

LOCKFILE=/var/lock/makewhatis.lock
[ -f $LOCKFILE ] && exit 0

trap "{ rm -f $LOCKFILE; exit 255; }" EXIT
touch $LOCKFILE
makewhatis -w
exit 0
