# libs/readline/profile.d/readline.rc

stty erase ^?

[ -f /etc/inputrc ] && export INPUTRC=/etc/inputrc

