# utils/findutils/profile.d/findutils.rc

# locate customized databases:
if [ -f ${HOME}/.locatedb ]; then
    export  LOCATE_PATH=/var/cache/locate/locatedb:${HOME}/.locatedb
fi

