# utils/ccache/profile.d/ccache.rc
# store the SA's ccache in a global location instead of /root/
if [ "$UID" == "0" ]; then
  export CCACHE_NOLINK=no
  export CCACHE_DIR=/var/cache/ccache
fi
