# utils/less/profile.d/less.rc
# these gadgets make less work nicer:
export LESS="-M -i -r"
export LESSOPEN="|lesspipe.sh %s"

