# crypto/openssh/init.d/sshd
#!/bin/bash
#
# ssh	Secure Shell Daemon
#
# chkconfig: 345 80 20
# description: SSH daemon
# processname: sshd
# config:      /etc/ssh/sshd_config

. /lib/lsb/init-functions

