# utils/which/profile.d/which.rc

# the *proper* way to use which:
which () {
    (alias; declare -f) | /usr/bin/which --tty-only --read-alias --read-functions --show-tilde --show-dot $@
}
export -f which

